# reference/python_full_loop_bundle_v15/dojo_bundle/systemd/mash-early-ssh.service
[Unit]
Description=MASH Early SSH (enable sshd + avahi) - before graphical login
DefaultDependencies=no
After=basic.target network-online.target
Wants=network-online.target
Before=sshd.service sddm.service display-manager.service graphical.target

[Service]
Type=oneshot
ExecStart=/usr/local/lib/mash/system/early-ssh.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
